# passthrough-watchdog.service
[Unit]
Description=Passthrough GPU Release Watchdog
Documentation=https://github.com/meeyao/one-click-passthrough
After=libvirtd.service
Wants=libvirtd.service

[Service]
Type=simple
Restart=on-failure
RestartSec=5
ExecStart=/usr/local/bin/passthrough-watchdog
StandardOutput=journal
StandardError=journal
SyslogIdentifier=passthrough-watchdog

[Install]
WantedBy=multi-user.target
