# passthrough-watchdog.service
[Unit]
Description=Passthrough GPU Release Watchdog
Documentation=https://github.com/meeyao/one-click-passthrough
After=libvirtd.service
Wants=libvirtd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/passthrough-status --json
ExecStop=/bin/true

[Install]
WantedBy=multi-user.target
